# cron.d/config-backup
0 * * * *  vyos /home/vyos/bin/config-backup > /home/vyos/bin/restore_config 2> /tmp/config-bakup.err
